# cron.d/clean_old_downloads
1 0 * * * root find /tmp/fm -mtime +1 -delete {} \;
